# fleet-configs/module2-microshift/etc/systemd/system/gateway-status-web.service
[Unit]
Description=MicroShift Gateway Status Web Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /etc/edge-config/gateway-status-web.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
