# fleet-configs/module2-microshift/edge-config-trigger.service
[Unit]
Description=Run MicroShift VRRP setup when fleet config arrives
ConditionPathExists=/etc/edge-config/setup-microshift.sh
ConditionPathExists=!/var/run/edge-config-microshift.done
After=network-online.target flightctl-agent.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /etc/edge-config/setup-microshift.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
